# hermes-skill-proxy.service
[Unit]
Description=Hermes Skill Proxy — auto-inject skill into system prompt
After=network.target hermes-gateway.service
Wants=hermes-gateway.service

[Service]
Type=simple
User=inexbot
WorkingDirectory=/home/inexbot
Environment=
  PROXY_PORT=8643
  HERMES_URL=http://localhost:8642
  HERMES_API_KEY=hermes-website-search
  SKILL_NAME=inexbot-knowledge-base

ExecStart=/usr/bin/python3 /home/inexbot/hermes-skill-proxy/hermes-skill-proxy.py
Restart=always
RestartSec=5

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes-skill-proxy

[Install]
WantedBy=multi-user.target
